# scheduler/com.theme_radar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>@@PLIST_LABEL@@</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@@PROJECT_ROOT@@/scripts/pipeline.sh</string>
    </array>

    <!-- 매일 06:00 실행 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- 표준 출력/에러 로그 -->
    <key>StandardOutPath</key>
    <string>@@PROJECT_ROOT@@/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>@@PROJECT_ROOT@@/logs/launchd_stderr.log</string>

    <!-- 환경변수 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@@HOME@@/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>@@HOME@@</string>
    </dict>

    <key>ExitTimeOut</key>
    <integer>3600</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
